# deploy/systemd/ink2md-purge.timer
[Unit]
Description=Run daily purge of Ink2MD output

[Timer]
OnCalendar=daily
AccuracySec=5m
Persistent=true
Unit=${SERVICE_NAME}-purge.service

[Install]
WantedBy=timers.target
