# deploy/systemd/ink2md.service
[Unit]
Description=Ink2MD processor service
Documentation=https://github.com/howardpchen/ink2md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${SERVICE_USER}
Group=${SERVICE_GROUP}
WorkingDirectory=${INSTALL_PREFIX}
EnvironmentFile=${ENV_FILE}
ExecStart=${INSTALL_PREFIX}/.venv/bin/ink2md --config ${CONFIG_PATH}
Restart=on-failure
RestartSec=10
RuntimeDirectory=ink2md
RuntimeDirectoryMode=0750
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
